# private_dot_config/systemd/user/crc-vsock.socket
[Unit]
Description=CRC vsock socket

[Socket]
ListenStream=vsock::1024
Service=crc-daemon.service

[Install]
WantedBy=default.target
